# docker/crontab
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

* * * * * root flock /var/lock/cron.lock --timeout 10 -c "/root/cron.sh"
