# docker/crontab
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

* * * * * root flock --timeout 10 /var/lock/cron.lock -c "/root/cron.sh"
